# etc/systemd/system/ofono.service.d/90-wait.conf
[Unit]
After=
After=droidian_boot_completed.service
